# gateway/systemd/voice-gateway.service
[Unit]
Description=Voice Gateway - AssemblyAI REST API Gateway
Documentation=https://github.com/prabhanshu11/omarchy-voice-typing
After=network.target

[Service]
Type=simple
# Kill any orphan process using port 8765 before starting
ExecStartPre=-/usr/bin/fuser -k 8765/tcp
ExecStart=%h/Programs/omarchy-voice-typing/gateway/voice-gateway
# Clean up on stop
ExecStopPost=-/usr/bin/fuser -k 8765/tcp
Restart=on-failure
RestartSec=5
WorkingDirectory=%h/Programs/omarchy-voice-typing/gateway

# Load API key from .env file
EnvironmentFile=%h/Programs/omarchy-voice-typing/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=voice-gateway

[Install]
WantedBy=default.target
